# config/openbad.service
[Unit]
Description=OpenBaD Agent Daemon
Documentation=https://github.com/bruceamoser/OpenBaD
After=openbad-broker.service network-online.target
Wants=network-online.target
Wants=openbad-broker.service

[Service]
Type=simple
User=openbad
Group=openbad
ExecStart=/usr/local/bin/openbad run --host localhost --port 1883
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s

# Environment
Environment=OPENBAD_CONFIG_DIR=/etc/openbad
WorkingDirectory=/var/lib/openbad

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/openbad /var/log/openbad
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

[Install]
WantedBy=multi-user.target
