# config/openbad-wui.service
[Unit]
Description=OpenBaD Web UI
Documentation=https://github.com/bruceamoser/OpenBaD
After=openbad.service network-online.target
Wants=network-online.target
Requires=openbad.service

[Service]
Type=simple
User=openbad
Group=openbad
ExecStart=/usr/local/bin/openbad wui --host 127.0.0.1 --port 9200 --mqtt-host localhost --mqtt-port 1883
Restart=on-failure
RestartSec=5s

Environment=OPENBAD_CONFIG_DIR=/etc/openbad
WorkingDirectory=/var/lib/openbad

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/openbad /var/log/openbad /etc/openbad
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

LimitNOFILE=65536
MemoryMax=256M

[Install]
WantedBy=multi-user.target